# device-agent-linux/device-agent-linux.service
[Unit]
Description=Device Agent for Linux
After=network.target

[Service]
ExecStartPre=/bin/bash -c 'sha256sum -c /usr/local/bin/device-agent-linux.sha256 || cp /backup/device-agent-linux /usr/local/bin/device-agent-linux'
ExecStart=/usr/local/bin/device-agent-linux
Restart=always
RestartSec=3
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target